# systemd/plumb-checks.service
[Unit]
Description=plumb: run deployment checks with AI escalation

[Service]
Type=oneshot
# APP is the deployment directory containing its pinned local plumb dependency.
Environment=PLUMB_CONFIG=%h/infra/APP/plumb.config
ExecStart=%h/infra/APP/node_modules/.bin/plumb run
# Alternative (unpinned network resolution): npx @kontourai/plumb run
