# systemd/plumb-checks.timer
[Unit]
Description=Daily plumb checks

[Timer]
OnCalendar=*-*-* 07:15:00
Persistent=true

[Install]
WantedBy=timers.target
